# setup/rules/99-automount_storage.rules
ACTION=="add", KERNEL=="sd?", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/home/pi/Repositories/cbs-ntcore/pi_monitor/scripts/automount_script.sh"
